# ozon-bi.service
[Unit]
Description=OZON Business Intelligence Bot
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/ozon-bi
Environment="PATH=/opt/ozon-bi/venv/bin"
EnvironmentFile=/opt/ozon-bi/.env
ExecStart=/opt/ozon-bi/venv/bin/python -m src.main

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=5min
StartLimitBurst=3

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ozon-bi

[Install]
WantedBy=multi-user.target
